# doc/systemd/ncontrol.target
[Unit]
Description=Devctrl Web Service and helpers
After=network.target auditd.service

[Install]
WantedBy=multi-user.target